# scheduler.crontab
20 12 * * * root cd /app && /usr/bin/node src/monitor.js && /usr/bin/node src/orchestrator.js >> /app/state/scheduler.log 2>&1
